# deploy/roomly-backend.service
[Unit]
Description=Roomly backend (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/roomly/backend
EnvironmentFile=/opt/roomly/.env
ExecStart=/opt/roomly/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
